# agents/options-trading/ai.openclaw.options-trading.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.options-trading.watchdog</string>
    <key>Comment</key>
    <string>Options Trading Agent Watchdog - monitors gateway, trading agent, and system health</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/sigbotti/.openclaw/workspace/agents/options-trading/watchdog.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>120</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/openclaw/watchdog-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/openclaw/watchdog-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/sigbotti</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>/Users/sigbotti/.openclaw/workspace/agents/options-trading</string>
    </dict>
</dict>
</plist>
